# config/crontab/chopsticks.crontab
*/2 * * * * /home/crawl/invi.sible.link/config/crontab/chopsticks
