# config/crontab/chopsticks.crontab
*/2 * * * * ISL="invi.sible.link" invi.sible.link/config/crontab/chopsticks
