# config/crontab/chopsticks.crontab
*/2 * * * * invi.sible.link/config/crontab/chopsticks
