# config/crontab/root.makesuremongo
#!/bin/sh
date +%c >> /tmp/makesuremongod
ps ax | grep mongod >> /tmp/makesuremongod
service mongod start

