# config/crontab/vigile.crontab
0 5 * * * /home/crawl/invi.sible.link/config/crontab/vigile-direction
