# config/crontab/vigile.crontab
0 5 * * * invi.sible.link/config/crontab/vigile-direction
