# bin/crontab/chopsticks.crontab
*/2 * * * * bin/ISL-scheduled --task chopsticks
