# bin/crontab/stats.crontab
*/10 * * * * bin/ISL-scheduled --task statistics
