# bin/crontab/vigile.crontab
0 5 * * * bin/ISL-scheduled --task vigile
